# home/.config/systemd/user/borgmatic.timer
[Unit]
Description=Run borgmatic backup

# Every day at 8pm
[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
