# home/.config/systemd/user/resticprofile-backup@profile-default.timer
[Unit]
Description=backup timer for profile default

[Timer]
OnCalendar=20:00
Persistent=true

[Install]
WantedBy=timers.target
