# home/.config/systemd/user/suspend.target
[Unit]
Description=This is started when the system is suspended
StopWhenUnneeded=yes
